# packages/thoradia/addons/webgrabplus/source/system.d/service.webgrabplus.service
[Unit]
Description=WebGrab+Plus
Documentation=http://www.webgrabplus.com/
After=kodi.service tools.mono.service
Requires=tools.mono.service

[Service]
ExecStart=/bin/sh /storage/.kodi/addons/service.webgrabplus/bin/webgrabplus.start
ExecStop=/usr/bin/systemctl stop webgrabplus.timer
ExecStop=/usr/bin/systemctl stop webgrabplus.service
RemainAfterExit=yes
Type=oneshot

[Install]
WantedBy=tools.mono.service
