# packages/thoradia/addons/nordvpn/source/system.d/service.nordvpn.service
Unit]
Description=NordVPN Daemon
After=basic.target network-online.target
Wants=basic.target network-online.target

[Service]
SyslogIdentifier=%N
ExecStart=/bin/sh /storage/.kodi/addons/%N/bin/nordvpn.ctl
KillMode=process
NonBlocking=true
Restart=on-failure
RestartSec=5

[Install]
WantedBy=basic.target
