# packages/thoradia/addons/mono/source/system.d/tools.mono.service
[Unit]
Description=Mono installer

[Service]
Type=oneshot
ExecStart=/bin/sh /storage/.kodi/addons/tools.mono/bin/tools.mono.install
RemainAfterExit=true
Restart=no

[Install]
WantedBy=kodi.target
