# packages/thoradia/addons/medusa/source/system.d/service.medusa.service
[Unit]
Description=Medusa
After=multi-user.target network-online.target
Wants=multi-user.target network-online.target

[Service]
SyslogIdentifier=%N
ExecStart=/bin/sh /storage/.kodi/addons/service.medusa/bin/medusa.ctl
Restart=always

[Install]
WantedBy=multi-user.target
