# packages/thoradia/addons/nzbget/source/system.d/service.nzbget.service
[Unit]
Description=NZBGet
After=multi-user.target network-online.target
Wants=multi-user.target network-online.target

[Service]
SyslogIdentifier=%N
Type=forking
ExecStart=/bin/sh /storage/.kodi/addons/service.nzbget/bin/nzbget.ctl -D
ExecStop=/bin/sh /storage/.kodi/addons/service.nzbget/bin/nzbget.ctl -Q

[Install]
WantedBy=multi-user.target
