# packages/thoradia/addons/rtorrent/source/system.d/service.rtorrent.service
[Unit]
Description=rTorrent
After=multi-user.target network-online.target
Wants=multi-user.target network-online.target

[Service]
SyslogIdentifier=%N
Type=forking
ExecStart=/bin/sh /storage/.kodi/addons/service.rtorrent/bin/rtorrent.screen start
ExecStop=/bin/sh /storage/.kodi/addons/service.rtorrent/bin/rtorrent.screen stop
Restart=always
KillMode=none

[Install]
WantedBy=multi-user.target
