# packages/thoradia/addons/deluge/source/system.d/service.deluge-web.service
[Unit]
Description=Deluge web user interface

[Service]
SyslogIdentifier=%N
ExecStart=/bin/sh /storage/.kodi/addons/service.deluge/bin/deluge-web.ctl -d

[Install]
WantedBy=service.deluge.service
