# packages/thoradia/addons/deluge/source/system.d/service.deluged.service
[Unit]
Description=Deluge daemon

[Service]
SyslogIdentifier=%N
ExecStart=/bin/sh /storage/.kodi/addons/service.deluge/bin/deluged.ctl -d

[Install]
WantedBy=service.deluge.service
