# packages/thoradia/addons/jackett/source/system.d/service.jackett.service
[Unit]
Description=Jackett
Documentation=https://github.com/Jackett/Jackett
After=kodi.service tools.mono.service
Requires=tools.mono.service
Wants=kodi.service

[Service]
ExecStart=/bin/sh /storage/.kodi/addons/service.jackett/bin/jackett.start
Restart=always

[Install]
WantedBy=kodi.service
