# packages/thoradia/addons/jackett/source/system.d/service.jackett.service
[Unit]
Description=Jackett
After=multi-user.target network-online.target tools.mono.service
Wants=multi-user.target network-online.target tools.mono.service

[Service]
SyslogIdentifier=%N
ExecStart=/bin/sh /storage/.kodi/addons/service.jackett/bin/jackett.ctl --NoRestart
Restart=always
MemoryHigh=125M

[Install]
WantedBy=multi-user.target
