# packages/thoradia/addons/transmission/source/system.d/service.transmission.service
[Unit]
Description=Transmission
After=multi-user.target network-online.target
Wants=multi-user.target network-online.target

[Service]
SyslogIdentifier=%N
Type=forking
ExecStart=/bin/sh /storage/.kodi/addons/service.transmission/bin/transmission.ctl
Restart=always

[Install]
WantedBy=multi-user.target
