# packages/thoradia/addons/lidarr/source/system.d/service.lidarr.service
[Unit]
Description=Lidarr
Documentation=http://lidarr.audio/
After=kodi.service tools.mono.service
Requires=tools.mono.service
Wants=kodi.service

[Service]
ExecStart=/bin/sh /storage/.kodi/addons/service.lidarr/bin/lidarr.start
Restart=always

[Install]
WantedBy=kodi.service
