# packages/thoradia/addons/bazarr/source/system.d/service.bazarr.service
[Unit]
Description=Bazarr
Documentation=https://github.com/morpheus65535/bazarr/wiki
After=multi-user.target network-online.target
Requires=multi-user.target network-online.target

[Service]
ExecStart=/bin/sh /storage/.kodi/addons/service.bazarr/bin/bazarr.start
Restart=always

[Install]
WantedBy=kodi.target
