# packages/thoradia/addons/couchpotato/source/system.d/service.couchpotato.service
[Unit]
Description=CouchPotato
Documentation=https://couchpota.to/
After=network.target
Requires=network.target

[Service]
ExecStart=/bin/sh /storage/.kodi/addons/service.couchpotato/bin/couchpotato.start
Restart=always

[Install]
WantedBy=kodi.target
