# packages/thoradia/addons/qbittorrent/source/system.d/service.qbittorrent.service
[Unit]
Description=qBittorrent service
Documentation=https://github.com/qbittorrent/qBittorrent/wiki
After=kodi.service
Wants=kodi.service

[Service]
Type=forking
ExecStart=/bin/sh /storage/.kodi/addons/service.qbittorrent/bin/qbittorrent.start
Restart=always

[Install]
WantedBy=kodi.service
