# packages/thoradia/addons/sonarr/source/system.d/service.sonarr.service
[Unit]
Description=Sonarr
After=multi-user.target network-online.target tools.mono.service
Wants=multi-user.target network-online.target tools.mono.service

[Service]
SyslogIdentifier=%N
ExecStart=/bin/sh /storage/.kodi/addons/service.sonarr/bin/sonarr.ctl -nobrowser
Restart=always
MemoryHigh=100M

[Install]
WantedBy=multi-user.target
