# packages/thoradia/addons/radarr/source/system.d/service.radarr.service
[Unit]
Description=Radarr
Documentation=https://github.com/Radarr/Radarr
After=kodi.service tools.mono.service
Requires=tools.mono.service
Wants=kodi.service

[Service]
ExecStart=/bin/sh /storage/.kodi/addons/service.radarr/bin/radarr.start
Restart=always
MemoryHigh=100M

[Install]
WantedBy=kodi.service
