# packages/thoradia/addons/radarr/source/system.d/service.radarr.service
[Unit]
Description=Radarr
After=multi-user.target network-online.target tools.mono.service
Wants=multi-user.target network-online.target tools.mono.service

[Service]
SyslogIdentifier=%N
ExecStart=/bin/sh /storage/.kodi/addons/service.radarr/bin/radarr.ctl -nobrowser
Restart=always
MemoryHigh=100M

[Install]
WantedBy=multi-user.target
